# rootfs/etc/s6/nginx/run
#!/bin/sh

exec /usr/sbin/nginx -p /usr/share/webapps
