# rootfs/etc/s6/crond/run
#!/bin/sh

if [[ -x /usr/bin/doas ]]; then
	exe=/usr/bin/doas
else
	exe=/usr/bin/sudo
fi
exec $exe /usr/sbin/crond -f
